# com.zeev.morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zeev.morning-brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/zeev.l/.claude/morning-brief/run.sh</string>
    </array>

    <!-- Sun(0) Mon(1) Tue(2) Wed(3) Thu(4) at 07:33 local — Israel work week -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>33</integer></dict>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>33</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>33</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>33</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>33</integer></dict>
    </array>

    <!-- if the mac was asleep at 07:33, run as soon as it wakes -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/zeev.l/.claude/morning-brief/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zeev.l/.claude/morning-brief/logs/launchd.err.log</string>
</dict>
</plist>
